# docker/etc/before_launch.ini
[program:before_launch]
command:sh /before_launch.sh
autorestart = false
startretries = 1
startsecs = 0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0